# provision/upstart/nomad.conf
description "Nomad by HashiCorp"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script

CONFIG_DIR=/etc/nomad.d
mkdir -p $CONFIG_DIR

# Grab ip for advertising purposes
ip=$(ifconfig eth1 | grep 'inet addr' | awk '{ print substr($2,6) }')
# Can't pass advertise flag so push into a config file, each node needs this to be its ip
cat > /etc/nomad.d/advertise.hcl <<END
advertise {
  rpc = "$ip:4647"
}
END

# region in consul = dc in nomad
region=$(hostname | cut -d'-' -f1)
# dc is a sub-region in nomad, so pretend we have different zones
dc=${region}-1

exec /usr/bin/nomad agent -region $region -dc $dc -config $CONFIG_DIR

end script

# Inspirations for this job file
# https://github.com/hashicorp/nomad/blob/master/demo/digitalocean/packer/nomad/upstart.nomad
